# tools/pve/encryption/zfs-load-all-keys-from-usb.service
[Unit]
Description=Load all ZFS encryption keys after USB key mount
DefaultDependencies=no
After=zfs-import.target mnt-_USB_PENDRIVE_KEY.mount
Requires=zfs-import.target mnt-_USB_PENDRIVE_KEY.mount
Before=zfs-mount.service
ConditionPathIsMountPoint=/mnt/_USB_PENDRIVE_KEY

[Service]
Type=oneshot
ExecStart=/usr/sbin/zfs load-key -a
TimeoutStartSec=120

[Install]
WantedBy=zfs.target

# Installation notes:
# - Ensure the USB device is mounted at /mnt/_USB_PENDRIVE_KEY before this service runs.
# - If using /etc/fstab, verify that systemd generates mnt-_USB_PENDRIVE_KEY.mount.
# - Copy this file to /etc/systemd/system/zfs-load-all-keys-from-usb.service
# - Then run: systemctl daemon-reload && systemctl enable --now zfs-load-all-keys-from-usb.service